# roles/rclone_sharepoint/templates/rclone-sharepoint.service.j2
[Unit]
Description=rclone mount {{ rclone_sharepoint_remote_name }} -> {{ rclone_sharepoint_mountpoint }}
After=network-online.target
Wants=network-online.target
AssertPathIsDirectory={{ rclone_sharepoint_mountpoint }}
ConditionPathExists={{ _rclone_sharepoint_config_path }}

[Service]
Type=notify
User={{ rclone_sharepoint_user }}
Group={{ rclone_sharepoint_user }}
# 起動時にマウント先を必ず用意して所有者を合わせる。
# "+" prefix を付けることで User=/Group= を無視して root で実行される。
ExecStartPre=+/bin/mkdir -p {{ rclone_sharepoint_mountpoint }}
ExecStartPre=+/bin/chown {{ rclone_sharepoint_user }}:{{ rclone_sharepoint_user }} {{ rclone_sharepoint_mountpoint }}
ExecStart=/usr/bin/rclone mount {{ rclone_sharepoint_remote_name }}: {{ rclone_sharepoint_mountpoint }} \
  --config {{ _rclone_sharepoint_config_path }} \
  --allow-other \
  --umask {{ rclone_sharepoint_umask }} \
  --dir-cache-time {{ rclone_sharepoint_dir_cache_time }} \
  --vfs-cache-mode {{ rclone_sharepoint_vfs_cache_mode }} \
  --log-level {{ rclone_sharepoint_log_level }}
ExecStop=/bin/fusermount3 -u {{ rclone_sharepoint_mountpoint }}
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
