# roles/llama/templates/llama-server.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=llama.cpp HTTP server (backend={{ llama_backend }})
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ansible_facts['env']['SUDO_USER'] }}
Group={{ ansible_facts['env']['SUDO_USER'] }}
WorkingDirectory={{ ansible_home }}/Llama/{{ llama_version }}
SupplementaryGroups=render video
{% if llama_backend == 'rocm' %}
Environment="HSA_OVERRIDE_GFX_VERSION={{ llama_server_hsa_override_gfx_version }}"
Environment="HSA_ENABLE_SDMA=0"
{% endif %}
ExecStart={{ ansible_home }}/Llama/{{ llama_version }}/build/bin/llama-server \
  --host {{ llama_server_host }} --port {{ llama_server_port }} \
  --model {{ llama_server_model_path }} \
  --n-gpu-layers {{ llama_server_n_gpu_layers }} \
  {{ llama_server_extra_args }}
{% if llama_backend in ['vulkan', 'rocm'] %}
DeviceAllow=/dev/dri/renderD128 rw
{% endif %}
{% if llama_backend == 'rocm' %}
DeviceAllow=/dev/kfd rw
{% endif %}
Restart=on-failure
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
